# docker/oracle/supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D

[program:oracle]
command=su - oracle -c 'lsnrctl start ' && su - oracle -c 'sqlplus / as sysdba' << EOF
startup 

EOF